# packages/container-lab/install/com.openai.codex-container-lab-reaper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openai.codex-container-lab-reaper</string>
  <key>ProgramArguments</key>
  <array>
    <string>__BUN_ABSOLUTE_PATH__</string>
    <string>__REAPER_ABSOLUTE_PATH__</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>StandardOutPath</key>
  <string>__STDOUT_LOG_ABSOLUTE_PATH__</string>
  <key>StandardErrorPath</key>
  <string>__STDERR_LOG_ABSOLUTE_PATH__</string>
</dict>
</plist>
